# scripts/lnmbot.service
[Unit]
Description=lnm-bot: BTC MA-cross strategy on LN Markets
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=james
Group=james
WorkingDirectory=/home/james/srv/tradingbot
EnvironmentFile=/etc/lnmbot/env
# Observe-only by default. Add --allow-orders only after the dedicated
# reconciliation smoke test; mainnet additionally requires --confirm-mainnet.
ExecStart=/home/james/srv/tradingbot/.venv/bin/python /home/james/srv/tradingbot/scripts/run_live.py --env /etc/lnmbot/env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lnmbot

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/lnmbot
Environment=STORAGE_DB_PATH=/var/lib/lnmbot/lnmarkets.sqlite

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
